# community/lxde/lxde-overlay/etc/systemd/system/multi-user.target.wants/myinit.service
[Unit]
Description=Init task

[Service]
ExecStart=/usr/bin/myinit

[Install]
WantedBy=multi-user.target
