# net-openrc/net-overlay/etc/local.d/local.stop
#!/bin/bash
#
# /etc/local.stop: Local shutdown script.
#

# Legacy sysvinit rc.local.shutdown support
if [ -e /etc/rc.local.shutdown ]; then
	# Execute it
	/etc/rc.local.shutdown
fi
