# net-openrc/net-overlay/etc/local.d/local.start
#!/bin/bash
#
# /etc/local.start: Local multi-user startup script.
#

# Legacy sysvinit rc.local support
if [ -e /etc/rc.local ]; then
	# Execute it
	/etc/rc.local
fi
