# xfce-openrc/xfce-overlay/etc/local.d/local.stop
#!/bin/bash
#
# /etc/local.d/local.stop: Local shutdown script.
#

# Legacy sysvinit rc.local.shutdown support
[ -x /etc/rc.local.shutdown ] && /etc/rc.local.shutdown
