# xfce-openrc/xfce-overlay/etc/local.d/local.start
#!/bin/bash
#
# /etc/local.d/local.start: Local multi-user startup script.
#

# Legacy sysvinit rc.local support
[ -x /etc/rc.local ] && /etc/rc.local
